# smf/manifests/binder-balancer.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2018, Joyent, Inc.
-->

<service_bundle type="manifest" name="binder-balancer">
    <service name="manta/application/binder-balancer" type="service"
      version="1">
        <create_default_instance enabled="true" />
        <single_instance />

        <dependency type="service" name="network" grouping="require_all"
          restart_on="error">
            <service_fmri value="svc:/milestone/network:default" />
        </dependency>

        <dependency type="service" name="filesystem" grouping="require_all"
          restart_on="error">
            <service_fmri value="svc:/system/filesystem/local" />
        </dependency>

        <dependency type="service" name="mdata" grouping="require_all"
          restart_on="none">
            <service_fmri value="svc:/smartdc/mdata:execute" />
        </dependency>

        <dependency type="service" name="config-agent" grouping="optional_all"
          restart_on="none">
            <service_fmri value="svc:/smartdc/application/config-agent" />
        </dependency>

        <exec_method type="method" name="start" timeout_seconds="10"
          exec="/opt/smartdc/binder/lib/balancer -s /var/run/binder/sockets
          &amp;">
            <method_context working_directory="/tmp">
                <method_credential user="nobody" group="nobody"
                  privileges="basic,net_privaddr" />
            </method_context>
        </exec_method>

        <exec_method type="method" name="refresh" exec=":kill"
          timeout_seconds="30" />

        <exec_method type="method" name="stop" exec=":kill"
          timeout_seconds="30" />

        <template>
            <common_name>
                <loctext xml:lang="C">Binder Loadbalancer</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
